# api/grid-paint-api.service
[Unit]
Description=Gridpaint

[Service]
ExecStart=exec-start.sh

[Install]
WantedBy=multi-user.target
