# api/grid-paint-api.service
[Unit]
Description=Gridpaint

[Service]
ExecStart=/home/eoin/grid-paint/api/exec-start.sh

[Install]
WantedBy=multi-user.target
